# sustainable-home-system/templates/tmux-main.service
[Unit]
Description=Start tmux server at login

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tmux-main-session

[Install]
WantedBy=default.target
